# config/staging/flask_project-staging.conf
[program:flask_project-staging]
command = gunicorn3 --certfile cert.pem --keyfile key.pem -b 0.0.0.0:8001 app:app
directory = /home/www/flask_project-staging/flask_project
user = olav
