# config/flask_staging.conf
[program:flask_staging]
command = gunicorn3 app:app -b localhost:8001
directory = /home/www/flask_staging/flask_staging
user = olav
