# config/production/flask_project.conf
[program:flask_project]
command = gunicorn3 app:app -b localhost:8000
directory = /home/www/flask_project/flask_project
user = olav
